# packaging/systemd/nl-router-receiver.service
[Unit]
Description=nl-router DICOM SCP (receiver)
Documentation=https://github.com/nl-router
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=nl-router
Group=nl-router

# Bootstrap config.
#
# /etc/nl-router/env is the shared file `nl-router init` generates from
# config.toml — it carries NL_ROUTER_SERVER_ID, NL_ROUTER_DATABASE_URL,
# and NL_ROUTER_LOG_LEVEL. /etc/nl-router/receiver.env is the optional
# per-service override file operators can hand-edit for receiver-only
# tuning (listen port, association timeout, etc.). Later entries
# override earlier ones, so the per-service file wins on conflict.
#
# Both files have leading `-` so the unit doesn't fail if the operator
# hasn't run `nl-router init` yet (the binary itself errors with a
# clear "required env var not set" message in that case).
EnvironmentFile=-/etc/nl-router/env
EnvironmentFile=-/etc/nl-router/receiver.env
ExecStart=/usr/bin/nl-receiver

Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s

# ---- Hardening -----------------------------------------------------------
# Drop every capability we don't need. The receiver binds DICOM ports
# (≥ 1024 by default; if operators choose < 1024 they can override this
# unit with AmbientCapabilities=CAP_NET_BIND_SERVICE).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
PrivateDevices=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Writable paths: landing zone + log dir. Receiver also reads its bootstrap
# config from /etc/nl-router which stays read-only.
ReadWritePaths=/var/lib/nl-router /var/log/nl-router

# Network: standard sockets only; no raw, no unusual address families.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Healthcheck-friendly: surface logs via journald in structured form.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
